# assignment04/10-kbd.rules
ACTION=="add", \
	SUBSYSTEM=="usb", \
	ATTRS{bInterfaceSubClass}=="01", \
	ATTRS{bInterfaceClass}=="03", \
	RUN+="/usr/sbin/insmod /kbd.ko"
